# example/null.socket.in
# SPDX-License-Identifier: GPL-2.0-or-later
#
# Copyright (C) 2026 Oracle.  All Rights Reserved.
# Author: Darrick J. Wong <djwong@kernel.org>
[Unit]
Description=Socket for null Service

[Socket]
ListenSequentialPacket=@FUSE_SERVICE_SOCKET_DIR_RAW@/null
Accept=yes
SocketMode=@FUSE_SERVICE_SOCKET_PERMS@
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
