# example/service_ll.socket.in
# SPDX-License-Identifier: GPL-2.0-or-later
#
# Copyright (C) 2026 Oracle.  All Rights Reserved.
# Author: Darrick J. Wong <djwong@kernel.org>
[Unit]
Description=Socket for service_ll Service

[Socket]
# Typically the socket path is /run/filesystems/service_ll
ListenSequentialPacket=@FUSE_SERVICE_SOCKET_DIR_RAW@/service_ll
Accept=yes
SocketMode=0660
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
